# runit/log/run
#!/bin/sh

mkdir -p /var/log/watchman/
exec svlogd -ttt /var/log/watchman/
